# alarmimg/usr/lib/systemd/system/oemcleanup.service
[Unit]
Description=Remove OEM user

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/userdel -r -f rebornos
ExecStart=/usr/bin/rm -rf /home/rebornos
ExecStart=/usr/bin/rm /etc/systemd/system/oemremove.service

[Install]
WantedBy=multi-user.target